# usr/share/rear/skel/default/lib/systemd/system/rear-boot-helper.service
[Unit]
Description=Rear boot script
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/etc/scripts/boot
