# usr/share/rear/skel/default/lib/systemd/system/rsyslog.service
[Unit]
Description=Relax-and-Recover run-syslog script
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/etc/scripts/run-syslog
