# usr/share/rear/skel/default/lib/systemd/system/run-system-setup.service
[Unit]
Description=Relax-and-Recover system setup script
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/etc/scripts/system-setup
