# usr/share/rear/skel/GALAXY11/etc/profile.d/galaxy11.sh
# set path to Commvault
PATH=$PATH:/opt/commvault/Base64
