# usr/share/rear/skel/default/usr/lib/systemd/system/sockets.target.wants/systemd-shutdownd.socket
#  This file is part of systemd.
#
# See systemd.special(7) for details

[Unit]
Description=Delayed Shutdown Socket
DefaultDependencies=no
Before=sockets.target

[Socket]
ListenDatagram=@/org/freedesktop/systemd1/shutdownd
