# usr/share/rear/skel/default/etc/inittab
# inittab for Relax-and-Recover
#
# Runlevel 3 is normal
#
id:3:initdefault:

bo::bootwait:/etc/scripts/boot
ss::bootwait:/etc/scripts/system-setup

ha:0:wait:/etc/scripts/do-shutdown halt
re:6:wait:/etc/scripts/do-shutdown reboot

sy:23:respawn:/etc/scripts/run-syslog
l3:3:wait:/etc/scripts/run-automatic-rear

t1:3:respawn:/sbin/mingetty --noclear tty1
t2:3:respawn:/sbin/mingetty --noclear tty2
t3:3:respawn:/sbin/mingetty --noclear tty3
t4:3:respawn:/sbin/mingetty --noclear tty4
t5:3:respawn:/sbin/mingetty --noclear tty5

