# usr/share/rear/skel/default/etc/init/start-syslog.conf
#
# This service starts the configured number of gettys.
start on startup
start on stopped rc RUNLEVEL=[2345]

task
script
	exec /etc/scripts/run-syslog
end script
