# usr/share/rear/skel/default/etc/init/start-serial.conf
### /etc/init/start-serial.conf
start on runlevel 2

task
script
    exec /etc/scripts/run-serial
end script
