# usr/share/rear/skel/default/etc/init/start-sshd.conf
### /etc/init/start-sshd.conf
start on runlevel 2

task
script
    exec /etc/scripts/run-sshd
end script
