# usr/share/rear/skel/NBKDC/etc/init/rcmd-executor.conf
description "NBK DC Remote Command Executor"
start on runlevel [2345]
stop on runlevel [016]

respawn
normal exit 0

exec /opt/NovaStor/DataCenter/rcmd-executor/rcmd-executor run
