# usr/share/rear/skel/default/usr/lib/systemd/system/getty.target
#  This file is part of systemd.
#

[Unit]
Description=Login Prompts
