# usr/share/rear/skel/default/usr/lib/systemd/system/halt.service
#  This file is part of systemd.
#

[Unit]
Description=Halt
DefaultDependencies=no
Requires=shutdown.target umount.target final.target
After=shutdown.target umount.target final.target

[Service]
Type=oneshot
ExecStart=/bin/systemctl --force halt
