# usr/share/rear/skel/default/usr/lib/systemd/system/rear-boot-helper.service
[Unit]
Description=Relax-and-Recover boot script
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/etc/scripts/boot
