# usr/share/rear/skel/default/usr/lib/systemd/system/rsyslog.service
[Unit]
Description=Relax-and-Recover run-syslog script

[Service]
Type=simple
ExecStart=/etc/scripts/run-syslog
