# usr/share/rear/skel/default/usr/lib/systemd/system/sysinit-opalpba.target
[Unit]
Description=TCG Opal pre-boot authentication (PBA)
Requires=sysinit-opalpba.service
AllowIsolate=yes

[Install]
Alias=sysinit-opalpba.target
