# usr/share/rear/skel/default/usr/lib/systemd/system/sysinit.service
[Unit]
Description=Initialize Rescue System
After=systemd-udevd.service

[Service]
Type=oneshot
ExecStart=/etc/scripts/system-setup
StandardInput=tty
RemainAfterExit=yes
