# usr/share/rear/skel/default/usr/lib/systemd/system/sysinit.service
[Unit]
Description=Initialize Rescue System

[Service]
Type=oneshot
ExecStart=/etc/scripts/system-setup
StandardInput=tty
